# deploy/vector/populate-vector-env.service
# Pulls AXIOM_PLATFORM_TOKEN + AXIOM_PLATFORM_DATASET (and metrics tokens)
# from Key Vault via the VM's managed identity and writes
# /etc/opensandbox/vector.env before vector.service starts.
#
# Reads OPENSANDBOX_AZURE_KEY_VAULT_NAME (and OPENCOMPUTER_CELL_ID /
# OPENSANDBOX_REGION for tagging) from whichever role-env file exists —
# worker.env on workers, server.env on the control plane. EnvironmentFile
# lines with a leading `-` are silently skipped if the file is absent, so
# this unit works on both roles.
#
# Azure first-boot race: cloud-init's cloud-final.service writes
# /etc/opensandbox/worker.env from a base64 payload baked by the control
# plane (internal/compute/azure.go). On Ubuntu Azure images, cloud-final
# is ordered After=multi-user.target — so worker.env doesn't exist when
# this unit runs (we're WantedBy=multi-user.target).
#
# Earlier attempts to handle this in-unit all failed:
#   #249: After=cloud-final.service → systemd cycle (vector silently dropped).
#   #254: exit 1 + Restart=on-failure → vector's restart-burst burnt the
#         StartLimit budget faster than RestartSec= could pace.
#   #256: internal 90s poll → multi-user blocked 90s, vector burnt its
#         own restart budget, cloud-final still hadn't run.
#   #257: internal 600s poll → boot deadlock (multi-user ↔ cloud-final),
#         every new Azure worker reaped at exactly 600s by pendingWorkerTTL.
#
# Current design (populate-vector-env.sh):
#   - If worker.env/server.env exists → populate synchronously, exit 0.
#   - If both missing → write stub vector.env, start
#     populate-vector-env-wait.service asynchronously, exit 0 in ~1s.
#   The wait unit polls in the background and re-runs the populator +
#   restarts vector once worker.env appears. Boot proceeds either way.

[Unit]
Description=Populate /etc/opensandbox/vector.env from Key Vault
After=network-online.target
Wants=network-online.target
Before=vector.service
DefaultDependencies=no

[Service]
Type=oneshot
EnvironmentFile=-/etc/opensandbox/worker.env
EnvironmentFile=-/etc/opensandbox/server.env
ExecStart=/usr/local/bin/populate-vector-env.sh
RemainAfterExit=yes
# Script exits 0 on every graceful path (missing role env → stub + kick
# waiter; missing IMDS / vault / secret → log + skip). Restart= would only
# fire on truly unexpected errors (set -e tripping), and even then the
# right move is to leave the unit in failed state for an operator to
# inspect rather than crash-loop.

[Install]
WantedBy=multi-user.target
